# nz7dev-gui.service
[Unit]
Description=NZ7DEV Mission Control GUI - Optimized Edition
Documentation=file:///README.md
After=network.target graphical-session.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/WD/dev/laika-local-rdp-nz7dev-gui
ExecStartPre=/bin/bash -c 'cd /mnt/WD/dev/laika-local-rdp-nz7dev-gui && chmod +x nz7dev && chmod +x launch_gui_optimized.sh'
ExecStart=/usr/bin/python3 -O -u /mnt/WD/dev/laika-local-rdp-nz7dev-gui/nz7dev_gui.py
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits
MemoryMax=1G
CPUQuota=75%

# Security settings (relaxed for development environment)
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# Optimized Environment
Environment=PYTHONPATH=/mnt/WD/dev/laika-local-rdp-nz7dev-gui
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=0
Environment=PYTHONOPTIMIZE=2
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=WERKZEUG_RUN_MAIN=true
Environment=GEVENT_SUPPORT=True
Environment=EVENTLET_NOPATCH=True
Environment=DISPLAY=:1
Environment=HOME=%h
Environment=PATH=/usr/bin:/usr/sbin:/bin:/sbin

# Environment variables for X11 display access - Updated for :1 display
Environment=XAUTHORITY=/home/nz7/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-1

# Logging
StandardOutput=append:/mnt/WD/dev/laika-local-rdp-nz7dev-gui/mission_control.log
StandardError=append:/mnt/WD/dev/laika-local-rdp-nz7dev-gui/mission_control.log
SyslogIdentifier=nz7dev-gui

# Ensure service runs in user session with display access
User=nz7
Group=nz7

[Install]
WantedBy=default.target 